# old_junk/Services/emolt_dat.service
[Unit]
After=mysql.target

[Service]
ExecStart=/etc/plumber/emolt_dat_refresh.sh

[Install]
WantedBy=default.target
